# core/docker/supervisord.conf
[supervisord]
user=root
nodaemon=true
logfile=/dev/stdout
logfile_maxbytes=0
loglevel=info

; ── 1. Virtual framebuffer ────────────────────────────────────────────────────
; Xvfb creates a virtual X11 display in RAM (:99). All GUI apps render to it.
; entrypoint.sh computes XVFB_RESOLUTION to include vertical offset headroom.
[program:xvfb]
command=Xvfb :99 -screen 0 %(ENV_XVFB_RESOLUTION)s -ac +extension GLX +render -noreset
autorestart=true
priority=100
startsecs=3
stopsignal=TERM
stopwaitsecs=10
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

; ── 2. Window manager ─────────────────────────────────────────────────────────
; Fluxbox provides window decorations and basic desktop management.
; Waits until the Xvfb display is actually accepting connections (xdpyinfo probe)
; instead of using an arbitrary sleep, eliminating timing-dependent failures.
[program:fluxbox]
command=bash -c "until DISPLAY=:99 xdpyinfo >/dev/null 2>&1; do sleep 0.5; done && DISPLAY=:99 fluxbox"
autorestart=true
priority=150
startsecs=3
stopsignal=TERM
stopwaitsecs=10
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

; ── 3. VNC password setup (runs once at container start) ─────────────────────
[program:vnc_setup]
command=bash -c "mkdir -p ~/.vnc && printf '%(ENV_VNC_PASSWORD)s' | vncpasswd -f > ~/.vnc/passwd && chmod 600 ~/.vnc/passwd"
autorestart=false
startsecs=0
priority=180
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

; ── 4. x11vnc — VNC server ───────────────────────────────────────────────────
; x11vnc reads the Xvfb framebuffer and speaks the VNC (RFB) protocol on :5901.
; Waits for the display to be ready before binding, preventing connection refusals
; that occurred with the previous sleep-based approach.
[program:x11vnc]
command=bash -c "until DISPLAY=:99 xdpyinfo >/dev/null 2>&1; do sleep 0.5; done && x11vnc -display :99 -forever -shared -rfbauth /root/.vnc/passwd -rfbport 5901 -noxrecord -noxfixes -noxdamage -noipv6"
autorestart=true
priority=200
startretries=10
startsecs=5
stopsignal=TERM
stopwaitsecs=10
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

; ── 5. noVNC — WebSocket proxy ───────────────────────────────────────────────
; noVNC proxies VNC (binary RFB) to WebSocket, and serves the HTML5 client.
; Waits until x11vnc is actually listening on :5901 before starting the proxy.
[program:novnc]
command=bash -c "until nc -z localhost 5901 >/dev/null 2>&1; do sleep 0.5; done && cd /opt/novnc && ./utils/novnc_proxy --vnc localhost:5901 --listen 0.0.0.0:8080 --web /opt/novnc"
autorestart=true
priority=300
startretries=5
startsecs=3
stopsignal=TERM
stopwaitsecs=10
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

; ── 6. Chrome ────────────────────────────────────────────────────────────────
; Profile data lives in /root/chrome-data (Docker named volume).
; Persists across container restarts and VM stop/start cycles.
; Lock files (SingletonCookie/Lock/Socket) are cleaned on each start;
; actual profile data (cookies, history, extensions) is never removed.
; To wipe the profile: docker compose down -v && docker compose up -d
;
; Waits for both Xvfb (xdpyinfo) and Fluxbox (pgrep) to be ready before
; launching, ensuring Chrome always has a stable window manager to render into.
[program:chrome]
command=bash -c "\
  mkdir -p /root/chrome-data/Default && \
  find /root/chrome-data -maxdepth 1 \
    \( -name SingletonCookie -o -name SingletonLock -o -name SingletonSocket \) \
    -delete 2>/dev/null || true && \
  until DISPLAY=:99 xdpyinfo >/dev/null 2>&1; do sleep 0.5; done && \
  until pgrep -x fluxbox >/dev/null 2>&1; do sleep 0.5; done && \
  DISPLAY=:99 /usr/bin/google-chrome-stable \
    --user-data-dir=/root/chrome-data \
    --profile-directory=Default \
    --window-position=0,%(ENV_WINDOW_OFFSET_Y)s \
    --window-size=%(ENV_RESOLUTION_WIDTH)s,%(ENV_RESOLUTION_HEIGHT)s \
    --no-sandbox \
    --disable-dev-shm-usage \
    --disable-gpu \
    --no-first-run \
    --no-default-browser-check \
    --disable-translate \
    --ignore-certificate-errors \
    --remote-debugging-port=9222 \
    --remote-debugging-address=0.0.0.0 \
    'about:blank'"
autostart=true
autorestart=true
priority=350
startretries=5
startsecs=10
stopsignal=TERM
stopwaitsecs=15
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
